# conf/supervisor-app.conf
[supervisord]
nodaemon=true

[program:app-uwsgi]
command = /usr/bin/uwsgi --ini /srv/graphsense-dashboard/graphsense-dashboard.ini --die-on-term

[program:nginx-app]
command = /usr/sbin/nginx
